# docker/scripts/supervisord.conf
[supervisord]
nodaemon=true

[program:init]
command=%(ENV_LIBSIMPLE_DIR)s/bin/initialize

[program:cron]
command=cron -f
